# luci-theme-Night/root/etc/uci-defaults/luci-theme-Night
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Night=/luci-static/Night
	set luci.main.mediaurlbase=/luci-static/Night
	commit luci
EOF
exit 0
